# ansible/roles/barman/templates/barman.shell.cron
#!/bin/bash
{% for key in backup_pg_cluster %}
{% if loop.first %}
barman backup {{ key }}
{% endif %}
{% endfor %}
